# root/etc/services.d/transmission/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/local/bin/transmission-daemon \
	-u $USER -v $PASSWD -p $WEBUI_PORT -P $PORT \
	-g /config -c /watch -f
